# deploy/linux/chronos-agent.service
[Unit]
Description=Chronos Agent — HTTP API for docker-compose on this host
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=chronos
Group=chronos
WorkingDirectory=/opt/chronos-agent
ExecStart=/usr/bin/env dotnet /opt/chronos-agent/Chronos.Agent.dll
Restart=always
RestartSec=5
TimeoutStopSec=30

# Optional: override URL, API key, paths, master registration (see chronos-agent.env.example)
Environment=ASPNETCORE_ENVIRONMENT=Production
EnvironmentFile=-/etc/chronos-agent/environment

# Compose projects and agent state
ReadWritePaths=/var/lib/chronos-agent

[Install]
WantedBy=multi-user.target
